# pi_testbench/services/reviision-pi.service
[Unit]
Description=ReViision Raspberry Pi Test Bench Service
Documentation=https://github.com/phloof/reviision
After=multi-user.target
After=network-online.target
Wants=network-online.target
Conflicts=shutdown.target
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/pi_testbench
Environment=PYTHONPATH=/home/admin/pi_testbench/src
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/admin/pi_testbench/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Use the wrapper script for reliable startup
ExecStart=/home/admin/pi_testbench/start_reviision.sh

# Graceful shutdown
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Restart configuration
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reviision-pi

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
MemoryMax=1G
TasksMax=200

# Environment file for configuration overrides
EnvironmentFile=-/home/admin/pi_testbench/.env

[Install]
WantedBy=multi-user.target
Alias=reviision.service 